# webinterface-onboot.service
[Unit]
Description=Enable the web interface on boot

[Service]
Type=simple
Restart=always
ExecStart=/home/root/.local/bin/webinterface-onboot

[Install]
WantedBy=multi-user.target